# ansible/roles/heads_tails_service/files/heads-tails.service
[Unit]
Description=heads-tails main python script
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/systemctl stop heads-tails-test
ExecStartPre=/usr/bin/pigpiod
ExecStart=/usr/bin/python3 /home/heads-tails/Python/heads-tails.py
ExecStopPost=/usr/bin/python3 /home/heads-tails/clearRegisters.py
ExecStopPost=/usr/bin/killall pigpiod
ExecStopPost=/usr/bin/sleep 5
Restart=always

[Install]
WantedBy=multi-user.target